# Docker/root/wechat-etc/supervisord.d/xybot-start.conf
[program:xybot-start]
command=xybot-start
autorestart=true